# devices/innotek-gmbh-virtualbox/files/etc/uci-defaults/10_disable_mqtt
#!/bin/sh
/etc/init.d/mosquitto disable
